# 97-hdmi.rules
ACTION=="change", SUBSYSTEM=="hdmi", RUN+="/bin/sh /home/pi/at.sh"